# dist/systemd/micro-bot.service
[Unit]
Description=Micro Bot micro-blogging IRC bot
Wants=network.target

[Service]
ExecStart=/usr/local/bin/micro-bot -config /root/.config/micro-bot/micro-bot.json

[Install]
WantedBy=default.target